# supervisor/osx.ini
[program:enomcms-reporter]
command=/usr/local/bin/python3 os_report.py
directory=/Users/james/code/emoncms-linux
autostart=true
autorestart=true
user=james
